# gentoo-init.d/ncpufreqd
#!/sbin/runscript

depend() {
	need localmount logger
}

start() {
	
	ebegin "Starting ncpufreqd"
	/usr/sbin/ncpufreqd
	eend $?
	
}

stop() {
	
	ebegin "Stopping ncpufreqd"
	killall ncpufreqd
	eend $?
	
}
