# gentoo-init.d/ncpufreqd
#!/sbin/runscript

depend() {
	need localmount logger
}

start() {

	ebegin "Starting ncpufreqd"
	start-stop-daemon --start --exec /usr/sbin/ncpufreqd
	eend $?

}

stop() {

	ebegin "Stopping ncpufreqd"
	start-stop-daemon --stop --exec /usr/sbin/ncpufreqd
	eend $?

}
